# backend.service
[Unit]
Description=API backend service

[Service]
ExecStart=%h/service/backend.sh
WorkingDirectory=%h/service

[Install]
WantedBy=default.target
